# port-forward.conf
description   "PIA VPN port forward to rtorrent"
author        "periwinklepreacher"

# The torrent-up event is generated by openvpn rtorrent-up.sh
start on torrent-up
stop on torrent-down

# These environment variables are assigned values by the event emitter
# (e.g., initctl emit torrent-up RT_LOGIN=/etc/openvpn/pia-login.ini
#        RT_CLIENTID=/etc/openvpn/pia-clientid.ini RT_IPV4=10.10.10.10
env RT_LOGIN
env RT_CLIENTID
env RT_IPV4
env RT_DEVICE

script
  log ( ) { logger "${UPSTART_JOB}[${UPSTART_EVENTS}]: $1"; }
  log "script started"

  USER=`head -n 1 $RT_LOGIN`
  PASSWD=`tail -n 1 $RT_LOGIN`
  CLIENTID=`head -n 1 $RT_CLIENTID`
  
  while [ 1 ]
  do
    URL='https://www.privateinternetaccess.com/vpninfo/port_forward_assignment'
    POSTDATA="user=${USER}&pass=${PASSWD}&client_id=${CLIENTID}&local_ip=${RT_IPV4}"
    RESPONSE=`wget -q --bind-address=${RT_IPV4} --post-data="${POSTDATA}" -O- ${URL}`
    if [ ! -z "${RESPONSE}" ]; then
      log "Response from PIA server is ${RESPONSE}"
    fi
    PORT=`echo ${RESPONSE} | grep -o "[0-9]*"`
    if [ -z "${PORT}" ]; then
      log "Failed to get port forward assignment for PIA VPN at ${RT_IPV4}."
    else
      log "Setting torrent port to ${PORT}"
      rtrpc --set set_port_range ${PORT}-${PORT}
      rtrpc --set set_port_random 0
    fi
    sleep 30m
  done
end script

post-stop script
  log ( ) { logger "${UPSTART_JOB}[${UPSTART_EVENTS}]: $1"; }
  log "post-stop script"
end script
